# agents/dev-machine-agent.service
[Unit]
Description=Device State Console Agent (Dev Machine)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=11

[Service]
Type=simple
WorkingDirectory=/home/lvziwang/OneDrive/文档/项目/设备状态控制台
EnvironmentFile=/home/lvziwang/OneDrive/文档/项目/设备状态控制台/agents/dev-machine-agent.env
ExecStart=/home/lvziwang/OneDrive/文档/项目/设备状态控制台/agents/run-dev-machine-agent.sh
Restart=on-failure
RestartSec=27

[Install]
WantedBy=default.target
