# agents/dev-machine-agent.service
[Unit]
Description=Device State Console Agent (Dev Machine)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=11

[Service]
Type=simple
WorkingDirectory=/opt/device-state-console-agent
EnvironmentFile=/etc/device-state-console-agent.env
ExecStart=/opt/device-state-console-agent/run-dev-machine-agent.sh
Restart=on-failure
RestartSec=27

[Install]
WantedBy=default.target
